# controlhub/pkg/cactuscore-controlhub.spec.nominal.filelist-finale
/usr/lib/systemd/system/controlhub.service
/etc/sysconfig/controlhub
